# systemd/user/calendar-sync.timer
[Unit]
Description=Sync CalDAV calendars every 30 minutes

[Timer]
OnCalendar=*:00/30
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
