# systemd/user/calendar-sync.service
[Unit]
Description=Sync CalDAV calendars and export to JSON
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/calendar-sync
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=default.target
