# systemd/user/dictation-server.service
[Unit]
Description=Remote dictation server (whisper.cpp over HTTPS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/dictation-server
Restart=on-failure
RestartSec=5
Environment=WHISPER_BIN=whisper-cli
Environment=DICTATION_FORMATTER_PROVIDER=openai
Environment=DICTATION_OPENAI_PROXY_URL=http://127.0.0.1:10531/v1
Environment=DICTATION_OPENAI_RP_MODEL=gpt-5.4
Environment=DICTATION_CLEANUP_MODEL=gpt-5.4-mini

[Install]
WantedBy=default.target
