# systemd/user/hermione.timer
[Unit]
Description=Hermione weekly audit — Sunday 09:00

[Timer]
OnCalendar=Sun 09:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
